# direct-pickup-checks/systemd/direct-pickup-monitor-summary.timer
[Unit]
Description=Daily direct-pickup-checks digest email

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
